# scripts/xpadcar_rpi.service
[Unit]
Description=xcar service

[Service]
ExecStart=/home/pi/xcar_rpi/scripts/connect_pad_and_run.sh

[Install]
WantedBy=multi-user.target

